# deployment-scripts/cmc-app.service
[Unit]
Description=CMC Cybersecurity Management Cockpit
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/cmc-app
ExecStart=/usr/bin/dotnet /var/www/cmc-app/CMC.Web.dll --urls http://localhost:5000
Restart=always
RestartSec=10
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
